# launchagents/com.kenscott.git-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kenscott.git-pull</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>
            cd ~/dotfiles
            before=$(git rev-parse HEAD)
            git pull --ff-only 2&gt;&amp;1
            after=$(git rev-parse HEAD)
            if [[ "$before" != "$after" ]]; then
                echo "dotfiles updated — relinking and updating vim plugins..."
                ~/dotfiles/update.sh 2&gt;&amp;1
            fi

            cd ~/dev/reminders-cli
            before=$(git rev-parse HEAD)
            git pull --ff-only 2&gt;&amp;1
            after=$(git rev-parse HEAD)
            if [[ "$before" != "$after" ]]; then
                echo "reminders-cli updated — rebuilding..."
                ~/dev/reminders-cli/reminders setup 2&gt;&amp;1
            fi
        </string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/com.kenscott.git-pull.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.kenscott.git-pull.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
