# launchagents/com.kenscott.backup-ai.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kenscott.backup-ai</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/ken.scott/bin/backup-ai-folder.py</string>
    </array>

    <!-- Run every 30 minutes -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <!-- Kill job cleanly if system goes to sleep mid-run -->
    <key>StandardAbandonProcessGroup</key>
    <true/>

    <!-- Script handles its own logging — discard stdout, capture unexpected stderr -->
    <key>StandardOutPath</key>
    <string>/dev/null</string>

    <key>StandardErrorPath</key>
    <string>/Users/ken.scott/logs/backup-ai-error.log</string>
</dict>
</plist>
